# debian/deepin-anything-server-service.postinst
#!/bin/bash
set -x
modprobe vfs_monitor || true
systemctl enable deepin-anything-server.service || true
systemctl start deepin-anything-server || true